# crontab
*/1 * * * * python /app/main.py >> /dev/stdout
